# dotfiles/scripts/crontab
* * * * * /home/cjber/scripts/odrive_sync.sh
* * * * * (sleep 30 ; /home/cjber/scripts/odrive_todo.sh )

