# config/systemd/user/restic-backup.timer
[Unit]
Description=Nightly restic backup to NAS

[Timer]
OnCalendar=*-*-* 03:00:00
RandomizedDelaySec=1800
Persistent=true

[Install]
WantedBy=timers.target
