# config/systemd/user/wt-cleanup.service
[Unit]
Description=Worktree orphan cleanup

[Service]
Type=oneshot
ExecStart=/home/cjber/scripts/wt_cleanup.sh
Nice=10
