# config/systemd/user/wt-cleanup.timer
[Unit]
Description=Daily worktree orphan cleanup

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target
