# dotfiles/config/autostart/mathpix-snipping-tool.desktop
[Desktop Entry]
Name=Mathpix Snipping Tool
Exec=mathpix-snipping-tool
Type=Application
Version=02.00.0047
