# scripts/crontab
0 * * * * /home/cjber/scripts/gdrive_sync.sh

