# scripts/crontab
* * * * * /home/cjber/scripts/gdrive_sync.sh 24h
5 * * * * /home/cjber/scripts/gdrive_sync.sh 5y

